# rpm-devel/tripwire/tripwire.cron.in
#!/bin/sh
HOST_NAME=`uname -n`
if [ ! -e @vardir@/lib/tripwire/${HOST_NAME}.twd ] ; then
	echo "****    Error: Tripwire database for ${HOST_NAME} not found.    ****"
	echo "**** Run "@sysconfdir@/tripwire/twinstall.sh" and/or "tripwire --init". ****"
else
        test -f @sysconfdir@/tripwire/tw.cfg &&  @sbindir@/tripwire --check
fi
